# config/hello_world_api.service.jinja
[Unit]
Description=hello world api service
After=network.target

[Service]
Restart=always
RestartSec=1
User={{ user }}
ExecStart={{ run_file }}

[Install]
WantedBy=multi-user.target